# update-r53@.timer
[Unit]
Description=Update AWS Route 53 record %i every 10 minutes

[Timer]
OnCalendar=*:0/30

[Install]
WantedBy=timers.target
